# client/systemd/host-ssh-certificate-renewal.timer
[Unit]
Description=Run SSH Host Key renewal

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
